# appliances/dhcp-ova/service/configure-guest.service
[Unit]
Description=DHCP Manager Guest Configuration

[Service]
Type=oneshot
ExecStart=/etc/photon/configure-guest.sh
ExecStop=
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
